# docker_compose/supervisor/supervisord.conf
[supervisord]
nodaemon=true


[program:gunicorn]
command=/usr/local/bin/gunicorn -c /code/gunicorn_conf.py --chdir /code/thuoclao thuoclao.wsgi:application --reload

autorestart=true
startretries=5